# battery-monitor/battery-monitor.conf
# myservice - myservice job file

description "It's a Python3 script which will notify user about charging, discharging and not charging state of the battery on Linux. It's using PyGtk3 for notification."
author "Maksudur Rahman Maateen <maateen@outlook.com>"

# Stanzas
#
# Stanzas control when and how a process is started and stopped
# See a list of stanzas here: http://upstart.ubuntu.com/wiki/Stanzas#respawn

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Essentially lets upstart know the process will detach itself to the background
expect fork

# Run before process
pre-start script
    exec python3 /usr/share/battery-monitor/battery-monitor.py
end script

# Start the process
exec myprocess
